# restdoc/package/etc/init/jon-deb-test.conf
# jon-deb-test.conf
start on filesystem
exec java -jar /opt/jinloes/jon-deb-test.jar --server.port=8181
